# hermes-bridge/systemd/hermes-bridge.service.example
[Unit]
Description=Hermes Bridge for Compass
After=network.target

[Service]
Type=simple

# IMPORTANT: bridge must run inside the Hermes Python environment so that
# `from run_agent import AIAgent` resolves. Adjust these three paths to match
# your install. Defaults below match a typical RackNerd/Ubuntu deploy where:
#   - hermes-bridge code lives in /opt/dash-repo/hermes-bridge
#   - Hermes Agent is installed at /root/.hermes/hermes-agent
WorkingDirectory=/opt/dash-repo/hermes-bridge
Environment=PYTHONPATH=/root/.hermes/hermes-agent
ExecStart=/root/.hermes/hermes-agent/venv/bin/python -m uvicorn app:app --host 127.0.0.1 --port 8787

# === Bridge-only settings ===
Environment=HOST=127.0.0.1
Environment=PORT=8787
# Set to a random secret if you want Compass to authenticate. Leave empty to
# allow loopback calls without a token.
Environment=HERMES_BRIDGE_TOKEN=
Environment=HERMES_BRIDGE_TIMEOUT_SECONDS=60
Environment=HERMES_BRIDGE_SKIP_MEMORY=false
Environment=HERMES_BRIDGE_SKIP_CONTEXT_FILES=false

# === Direct fallback (when Hermes Python path fails) ===
# Bridge first tries to call the Hermes AIAgent. If that fails (e.g. config.yaml
# misconfigured, missing key, returns "None"), bridge falls back to calling the
# OpenAI-compatible endpoint directly. Set the values below to match what you
# want bridge's fallback to hit. They DO NOT have to match config.yaml's primary,
# but they should usually mirror config.yaml's fallback_model so user experience
# stays consistent across paths.
Environment=HERMES_BRIDGE_ALLOW_DIRECT_FALLBACK=true
Environment=HERMES_BRIDGE_DIRECT_FOR_ATTACHMENTS=true

# Default fallback target: SiliconFlow + DeepSeek-V3.2 (matches the recommended
# fallback_model in ~/.hermes/config.yaml). Override per deployment as needed.
Environment=HERMES_BRIDGE_FALLBACK_BASE_URL=https://api.siliconflow.cn/v1
Environment=HERMES_BRIDGE_FALLBACK_MODEL=deepseek-ai/DeepSeek-V3.2

# === HERMES_AGENT_* overrides AIAgent constructor args ===
# Use these only when you can NOT modify ~/.hermes/config.yaml and need bridge
# to coerce AIAgent into a specific provider. Leave them out otherwise so
# AIAgent reads config.yaml normally — same as hermes-gateway does.
# Example (NVIDIA NIM as primary, matching config.yaml):
# Environment=HERMES_AGENT_PROVIDER=custom
# Environment=HERMES_AGENT_BASE_URL=https://integrate.api.nvidia.com/v1
# Environment=HERMES_AGENT_MODEL=moonshotai/kimi-k2-instruct

# === API keys: shared via runtime.env (NEVER inline here) ===
# hermes-gateway uses /root/.hermes/runtime.env. Putting keys ONLY there means
# both services read the same source. Specifically expected entries (do not
# commit this file to git):
#   NVIDIA_API_KEY=nvapi-...
#   SILICONFLOW_API_KEY=sk-...
#   HERMES_BRIDGE_FALLBACK_API_KEY=sk-...    (← bridge's direct fallback key,
#                                               typically same as SILICONFLOW_API_KEY)
EnvironmentFile=-/root/.hermes/runtime.env

# DO NOT add `Environment=*_API_KEY=...` lines here. Those would shadow
# runtime.env, only affect this service, and silently diverge from
# hermes-gateway, causing "Web works but Discord fails" split-brain bugs.

Restart=on-failure
RestartSec=3
MemoryMax=512M

[Install]
WantedBy=multi-user.target
